# src/main/resources/dockerfiles/labrador/centos7/supervisord.conf

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:postgresql]
command=/bin/bash -c "/usr/bin/postgres -D /var/lib/pgsql/data -p 5432"
user=postgres
autostart=true
autorestart=false
startretries=0
startsecs=0
priority=100

[program:populatedbdata]
user=root
autostart=true
autorestart=false
command=/bin/bash -c "( /bin/sleep 5; /home/deploy_db_data.sh; )"
startretries=0
startsecs=0
priority=999


[program:jetty]
user=root
command=java -jar /opt/jetty/start.jar
autostart=true
autorestart=false
startretries=0
startsecs=0
priority=300

